# conf/supervisor/ampel.conf
[program:ampel]
command=/usr/bin/python /home/pi/py_ws_xfd/main.py
autorestart=true
user=root
directory=/home/pi/py_ws_xfd
