# conf/supervisor/networkmanager.conf
[program:networkmanager]
command=/usr/bin/python /home/pi/py_ws_xfd/networkmanager.py
autorestart=true
user=root
directory=/home/pi/py_ws_xfd
